# Home/.config/autostart/mountMain.sh.desktop
[Desktop Entry]
Exec=/home/akshay/Documents/mountMain.sh
Icon=dialog-scripts
Name=mountMain.sh
Path=
Type=Application
X-KDE-AutostartScript=true
